# MTBSystem/dockerbase/supervisor/api-srv.conf
[program:api-srv]
command=/data/deploy/mtbsystem/src/api-srv/api-srv
process_name=api-srv
autorestart=true
redirect_stderr=true
stderr_logfile=/data/logs/api-srv.err.log
stdout_logfile=/data/logs/api-srv.out.log
user=root