# MTBSystem/dockerbase/supervisor/comment-srv.conf
[program:comment-srv]
command=/data/deploy/mtbsystem/src/comment-srv/comment-srv
process_name=comment-srv
autorestart=true
redirect_stderr=true
stderr_logfile=/data/logs/comment-srv.err.log
stdout_logfile=/data/logs/comment-srv.out.log
user=root